# monitor_test.service
[Unit]
Description=Мониторинг процесса test
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/process_monitor_test.sh
